# kaliber-server.service
[Unit]
Description=Kaliber eBook Server
After=network.target

[Service]
Type=idle
User=matthias
Group=matthias
WorkingDirectory=/home/matthias/devel/Go/src/github.com/mwat56/kaliber/
ExecStart=/home/matthias/devel/Go/src/github.com/mwat56/kaliber/bin/kaliber-linux-amd64 -listen=0 -log=access.read.mwat.de
Restart=on-failure

[Install]
WantedBy=multi-user.target
